# scripts/install/linux/dahua-adapter.service
[Unit]
Description=Dahua Adapter
After=network.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=__NODE_BIN__ dist/src/index.js
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
